# supervisord.conf
[supervisord]
pidfile     = /var/log/supervisord.pid
logfile     = /var/log/supervisord.log
childlogdir = /var/log/

logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=warn                                   ; info, debug, warn, trace

minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors

logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace

nodaemon=false                                  ; run supervisord as a daemon
user=root                                       ; default user

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[inet_http_server]
port=0.0.0.0:9001

[program:dashboard]
environment                     = NODE_ENV=production
command                         = /usr/local/bin/node ./frontend/express/app.js
numprocs                        = 1
numprocs_start                  = 0
directory                       = /opt/countly
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
stopasgroup                     = true
killasgroup                     = true
redirect_stderr                 = true
stdout_logfile                  = /var/log/countly/dashboard.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = /var/log/countly/dashboard.log
stdout_events_enabled           = false
loglevel                        = warn



[program:api]
environment                     = NODE_ENV=production
command                         = /usr/local/bin/node ./api/api.js
numprocs                        = 1
numprocs_start                  = 0
directory                       = /opt/countly
autostart                       = true
autorestart                     = true
startsecs                       = 10
startretries                    = 5
stopsignal                      = TERM
stopwaitsecs                    = 10
stopasgroup                     = true
killasgroup                     = true
redirect_stderr                 = true
stdout_logfile                  = /var/log/countly/api.log
stdout_logfile_maxbytes         = 1MB
stdout_logfile_backups          = 10
stdout_capture_maxbytes         = 1MB
stderr_logfile                  = /var/log/countly/api.log
stdout_events_enabled           = false
loglevel                        = warn


[group:countly]
programs=dashboard, api
